# python_scripts/crontab
*/1 * * * * python /scraper.py >> /logfile